# scripts/streamsnatcher.service
[Unit]
Description=StreamSnatcher - P2P File Transfer
After=network.target

[Service]
Type=simple
User=streamsnatcher
Group=streamsnatcher
WorkingDirectory=/home/streamsnatcher/htdocs/streamsnatcher.com
Environment="PATH=/home/streamsnatcher/htdocs/streamsnatcher.com/venv/bin:/usr/bin"
EnvironmentFile=/home/streamsnatcher/htdocs/streamsnatcher.com/.env.production
ExecStart=/home/streamsnatcher/htdocs/streamsnatcher.com/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 4
Restart=always
RestartSec=5

# Security hardening (compatible with container-based VPS)
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
